# ansible/flask/supervisord.conf
[supervisord]
nodaemon=true

[program:python]
command=python app.py
